# gdbus_p2p-1.0/data/gdbusp2p.service.in
[Unit]
Description=gdbus p2p service

[Service]
#Type=dbus
#BusName=org.freedesktop.Gdbusp2p
ExecStart=@gdbusp2pprivdir@/gdbusp2pd
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target